# debian/droid4-wlanconfig.init
#!/sbin/openrc-run

description="Maserati Wifi calibration"

depend() {
	need udev
}

start() {
	ebegin "Calibrating Wifi"
	if /usr/bin/maserati-calibrate; then
		rc-update del droid4-wlanconfig
	fi
	eend $?
}

stop() {
	return 0
}
